# deploy/kkj-jpdir.timer
[Unit]
Description=Sync the Japanese machine-readable resource directory daily

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
