# deploy/kkj-poll.service
[Unit]
Description=kkj-watch slow lane (source documents + extraction + backup)

[Service]
Type=oneshot
User=kkj
WorkingDirectory=/opt/kkj-watch
EnvironmentFile=/etc/kkj-watch.env
ExecStart=/usr/bin/python3 -m kkj.pipeline poll-docs
ExecStart=/usr/bin/python3 -m kkj.pipeline extract 300
ExecStart=/usr/bin/python3 -m kkj.pipeline backup
