# deploy/kkj-poll-fast.timer
[Unit]
Description=Run kkj-watch fast lane hourly

[Timer]
OnCalendar=*-*-* *:15:00
Persistent=true

[Install]
WantedBy=timers.target
