# deploy/kkj-attest.service
[Unit]
Description=kkj-watch daily signed hash-chain root (x402 Trust Index attestation)

[Service]
Type=oneshot
User=kkj
WorkingDirectory=/opt/kkj-watch
EnvironmentFile=/etc/kkj-watch.env
ExecStart=/usr/bin/python3 -m kkj.attest root
