# deploy/kkj-poll-fast.service
[Unit]
Description=kkj-watch fast lane (portal poll + correction linking + semantic analysis)

[Service]
Type=oneshot
User=kkj
WorkingDirectory=/opt/kkj-watch
EnvironmentFile=/etc/kkj-watch.env
ExecStart=/usr/bin/python3 -m kkj.pipeline poll
ExecStart=/usr/bin/python3 -m kkj.pipeline link 50
