# deploy/kkj-poll.service
[Unit]
Description=kkj-watch polling (cases + source documents + extraction)

[Service]
Type=oneshot
User=kkj
WorkingDirectory=/opt/kkj-watch
EnvironmentFile=/etc/kkj-watch.env
ExecStart=/usr/bin/python3 -m kkj.pipeline poll
ExecStart=/usr/bin/python3 -m kkj.pipeline poll-docs
ExecStart=/usr/bin/python3 -m kkj.pipeline extract 20
