# deploy/kkj-poll.timer
[Unit]
Description=Run kkj-watch polling every 3 hours

[Timer]
OnCalendar=00/3:10
Persistent=true

[Install]
WantedBy=timers.target
